# scripts/common/snapmulti-mpd-update.service
[Unit]
Description=Trigger MPD incremental rescan (NFS/SMB libraries only)
# Need docker up so `docker exec mpd` works, and snapmulti-server.service
# so the mpd container is actually started.
After=docker.service snapmulti-server.service
# Don't pile up on transient failures (mpd container restarting, NFS
# transient). 3 attempts inside a 10-min window is generous; the next
# nightly firing handles persistent failures cleanly.
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=oneshot
# 35 min: covers the script's internal 30-min poll deadline plus a
# small margin for Docker exec startup latency.
TimeoutStartSec=2100
ExecStart=/usr/local/bin/mpd-nfs-update
