# scripts/common/snapmulti-state-backup.service
[Unit]
Description=Backup snapMULTI runtime state to boot partition
After=local-fs.target
ConditionPathExists=/opt/snapmulti

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-snapmulti-state
