# scripts/common/snapmulti-diagnostics.timer
[Unit]
Description=Periodic snapMULTI diagnostic log capture

[Timer]
# Save diagnostics every 15 minutes and on boot.
# Retention capped at 12 snapshots (3 h of context) inside save-diagnostics.
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
