# scripts/common/snapmulti-state-backup.timer
[Unit]
Description=Periodic snapshot of snapserver + myMPD state (primary mechanism for server.json)
Documentation=man:systemd.timer(5)

[Timer]
# Primary backup mechanism for server.json (the .path unit no longer
# watches it — snapserver rewrites it every ~3 s with refreshed
# lastSeen heartbeats which would cause a backup loop). Also acts as a
# safety net for myMPD subdirectory writes that .path units don't fire
# on (systemd path units are not recursive). backup-snapmulti-state.sh
# is idempotent and short-circuits on canonical equality (server.json
# minus lastSeen), so a no-op tick costs only a JSON parse + hash —
# zero FAT32 writes.
OnBootSec=5min
OnUnitActiveSec=5min
Unit=snapmulti-state-backup.service
Persistent=true

[Install]
WantedBy=timers.target
