# scripts/common/snapmulti-backup.timer
[Unit]
Description=Periodic MPD database backup to boot partition

[Timer]
# Backup daily at 4 AM and on boot (after containers are up). Runtime
# snapserver/myMPD state is event-backed by snapmulti-state-backup.path.
OnBootSec=5min
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
