# scripts/common/snapmulti-mpd-update.timer
[Unit]
Description=Nightly MPD incremental rescan for NFS/SMB libraries

[Timer]
# 03:00 daily — 1 h before snapmulti-backup.timer at 04:00 so the
# next morning's backup picks up freshly-scanned tracks. Persistent=true
# catches missed runs after shutdown / suspend / power outage.
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
