# client/common/systemd/snapclient-display.service
[Unit]
Description=Snapclient display detection (HDMI profile for .env)
Before=snapclient.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/opt/snapclient/scripts/display-detect.sh
WorkingDirectory=/opt/snapclient
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
